# debian/mydaemon.service
[Unit]
Description=My Python Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m mydaemon.main
Restart=always
RestartSec=10
User=nobody
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=mydaemon

[Install]
WantedBy=multi-user.target
    