# debian/prerm
#!/bin/sh
set -e

if [ "$1" = "remove" ]; then
    # Stop the service before removal
    systemctl stop mydaemon.service || true
    systemctl disable mydaemon.service || true
fi

#DEBHELPER#
